# packaging/debian/postinst
#!/bin/bash
#  * Project: OpenAuto
#  * This file is part of openauto project.
#  * Copyright (C) 2025 OpenCarDev Team
#  *
#  *  openauto is free software: you can redistribute it and/or modify
#  *  it under the terms of the GNU General Public License as published by
#  *  the Free Software Foundation; either version 3 of the License, or
#  *  (at your option) any later version.
#  *
#  *  openauto is distributed in the hope that it will be useful,
#  *  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  *  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  *  GNU General Public License for more details.
#  *
#  *  You should have received a copy of the GNU General Public License
#  *  along with openauto. If not, see <http://www.gnu.org/licenses/>.

set -e

case "$1" in
    configure)
        # Reload udev rules
        if [ -x /bin/udevadm ]; then
            /bin/udevadm control --reload-rules || true
            /bin/udevadm trigger || true
        fi

        # Reload systemd daemon
        if [ -x /bin/systemctl ]; then
            /bin/systemctl daemon-reload || true
        fi

        # Stop and mask getty@tty2 to prevent interference with linuxfb GUI
        if [ -x /bin/systemctl ]; then
            echo "Configuring tty2 for OpenAuto linuxfb output..."
            /bin/systemctl stop getty@tty2.service 2>/dev/null || true
            /bin/systemctl mask getty@tty2.service 2>/dev/null || true
        fi

        # Create ld.so.conf.d entry for custom libraries
        if [ -d /opt/openauto/lib ]; then
            echo "/opt/openauto/lib" > /etc/ld.so.conf.d/openauto.conf
            ldconfig || true
        fi

        # Disable console blanking for framebuffer
        if [ -f /sys/class/graphics/fbcon/cursor_blink ]; then
            echo 0 > /sys/class/graphics/fbcon/cursor_blink || true
        fi

        # Enable Bluetooth service (main app uses bash_profile autologin)
        if [ -x /bin/systemctl ]; then
            /bin/systemctl enable openauto-btservice.service || true
        fi

        echo ""
        echo "OpenAuto (RK3229/Armbian) installed successfully!"
        echo ""
        echo "Configuration:"
        echo "  - GUI: Qt EGLFS on DRM/KMS"
        echo "  - Video: FFmpeg DRM hwaccel (zero-copy)"
        echo "  - Audio: RTAudio ALSA (direct, no PulseAudio)"
        echo "  - Libraries: /opt/openauto/lib"
        echo ""
        echo "OpenAuto will start automatically on login via ~/.bash_profile"
        echo "To start Bluetooth service manually:"
        echo "  sudo systemctl start openauto-btservice"
        echo ""
        echo "Note: getty@tty2 has been masked to prevent interference with the GUI."
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
