# packaging/systemd/openauto.service
#  * Project: OpenAuto (RK3229/Armbian Custom Build)
#  * This file is part of openauto project.
#  * Copyright (C) 2025 OpenCarDev Team
#  *
#  *  openauto is free software: you can redistribute it and/or modify
#  *  it under the terms of the GNU General Public License as published by
#  *  the Free Software Foundation; either version 3 of the License, or
#  *  (at your option) any later version.
#  *
#  *  openauto is distributed in the hope that it will be useful,
#  *  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  *  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  *  GNU General Public License for more details.
#  *
#  *  You should have received a copy of the GNU General Public License
#  *  along with openauto. If not, see <http://www.gnu.org/licenses/>.

[Unit]
Description=OpenAuto Android Auto Head Unit (RK3229/Armbian)
After=network.target sound.target bluetooth.target local-fs.target
# Prevent getty from interfering with linuxfb on tty2
Conflicts=getty@tty2.service
After=getty@tty2.service

[Service]
Type=simple
# Use launcher script that configures linuxfb, kmssink, and RTAudio
ExecStart=/opt/openauto/openauto.sh
Restart=on-failure
RestartSec=5s
User=root

# Run on tty2 for linuxfb output
TTYPath=/dev/tty2
StandardInput=tty
StandardOutput=tty
TTYVTDisallocate=yes

# Qt linuxfb configuration
Environment="QT_QPA_PLATFORM=linuxfb:fb=/dev/fb0:tty=/dev/tty2"
Environment="QT_QPA_FB_DRM=1"
Environment="QT_QPA_FB_NO_LIBINPUT=0"

# Library path for custom dependencies (aasdk, aap_protobuf, etc.)
Environment="LD_LIBRARY_PATH=/opt/openauto/lib:/usr/local/lib"

# GStreamer kmssink configuration for RK3229
Environment="GST_DEBUG=2"
Environment="GST_PLUGIN_PATH=/usr/lib/arm-linux-gnueabihf/gstreamer-1.0"

# RTAudio ALSA configuration (bypass PulseAudio)
Environment="RTAUDIO_API=alsa"
Environment="PULSE_SERVER=none"

# Runtime directory
Environment="XDG_RUNTIME_DIR=/run/user/0"

# Ensure audio/video/input access
SupplementaryGroups=audio video plugdev input tty

# Prepare tty2 before starting
ExecStartPre=/bin/chvt 2
ExecStartPre=-/bin/setterm --cursor off --clear all </dev/tty2 >/dev/tty2
ExecStartPre=-/bin/sh -c 'echo 0 > /sys/class/graphics/fbcon/cursor_blink'

[Install]
WantedBy=multi-user.target
